# z Plexmedia install/supervisor/conf.d/openvpn.conf
[program:openvpn]
command=/usr/sbin/openvpn --client --dev tun0 --config /etc/openvpn/client.conf