# z Plexmedia install/supervisor/conf.d/syslog.conf
[program:syslog]
command=/usr/sbin/rsyslogd -n